# bot_project/hls-explorer.service
[Unit]
Description=Telegram HLS Live Explorer Bot
After=network.target

[Service]
User=root
WorkingDirectory=/opt/hls-explorer
ExecStart=/opt/hls-explorer/venv/bin/python -m app.main
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
